# etc/cron.d/easynas.cron
0 */6 * * *  root sleep ${RANDOM:0:2}m ;/easynas/startup/check_update.pl
